# gitorious/templates/default/git-ultrasphinx.erb
#! /bin/sh
### BEGIN INIT INFO
# Provides:          Gitorious-Ultrasphinx
# Required-Start:    mysql
# Required-Stop:     
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Ultrasphinx daemon
# Description:       Starts the Ultrasphinx daemons needed by Gitorious
### END INIT INFO

# Generated by Chef for <%= node[:fqdn] %>
# Local modifications will be overwritten.

RETVAL=0
START_CMD="  cd <%= @deploy_path %> && bundle exec rake ultrasphinx:daemon:start   RAILS_ENV=production"
STOP_CMD="   cd <%= @deploy_path %> && bundle exec rake ultrasphinx:daemon:stop    RAILS_ENV=production"
RESTART_CMD="cd <%= @deploy_path %> && bundle exec rake ultrasphinx:daemon:restart RAILS_ENV=production"
STATUS_CMD=" cd <%= @deploy_path %> && bundle exec rake ultrasphinx:daemon:status  RAILS_ENV=production" 
LOCK_FILE=/var/lock/git-ultrasphinx
PID_FILE=<%= @deploy_path %>/db/sphinx/log/searchd.pid

case "$1" in
  start)
    /bin/su git -c "$START_CMD"
    ;;
  stop)
    /bin/su git -c "$STOP_CMD"
    ;;
  status)
    /bin/su git -c "$STATUS_CMD"
    ;;
        restart)
    /bin/su git -c "$RESTART_CMD"
                ;;
  *)
    echo $"Usage: $0 {start|stop|restart|status}"
    RETVAL=1
esac
exit $RETVAL
